# systemd/dmc-ai-mobility.service
[Unit]
Description=DMC AI Mobility robot node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/dmc_ai_mobility
Environment=PYTHONPATH=/opt/dmc_ai_mobility/src
ExecStart=/bin/sh -lc 'if command -v libcamerify >/dev/null 2>&1; then exec libcamerify /usr/bin/python3 -m dmc_ai_mobility.app.cli robot --config /opt/dmc_ai_mobility/config.toml; else exec /usr/bin/python3 -m dmc_ai_mobility.app.cli robot --config /opt/dmc_ai_mobility/config.toml; fi'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
