# systemd/dmc-ai-mobility.service
[Unit]
Description=DMC AI Mobility robot node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/fooping/dmc_ai_mobility
Environment=VIRTUAL_ENV=/home/fooping/env
Environment=PATH=/home/fooping/env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/fooping/dmc_ai_mobility/src
ExecStart=/home/fooping/dmc_ai_mobility/scripts/run_robot.sh /home/fooping/dmc_ai_mobility/config.toml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
