# systemd/dmc-ai-mobility.service
[Unit]
Description=DMC AI Mobility robot node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/dmc_ai_mobility
Environment=PYTHONPATH=/opt/dmc_ai_mobility/src
ExecStart=/opt/dmc_ai_mobility/scripts/run_robot.sh /opt/dmc_ai_mobility/config.toml
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
